# deploy/python_ec2/music-app-web.service
[Unit]
Description=Oracle music recommendation FastAPI app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/music_app/app
EnvironmentFile=/home/ubuntu/music_app/app/.env
ExecStart=/home/ubuntu/music_app/.venv/bin/uvicorn 04_web_preview.app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
